# scripts/deluge/template.deluged.conf
# deluged - Deluge daemon
#
# The daemon component of Deluge BitTorrent client. Deluge UI clients
# connect to this daemon via DelugeRPC protocol.

description "Deluge daemon for <user>"
author "Seedbox automation script"

start on filesystem and static-network-up
stop on runlevel [016]

respawn
respawn limit 5 30

env uid=<user>
env gid=<user>
env umask=022

exec start-stop-daemon -S -c $uid:$gid -k $umask -x /usr/bin/deluged -- -d
